# services.d/postgis-proxy/run
#!/usr/bin/with-contenv bash
exec 2>&1
cd /tmp
exec socat TCP-LISTEN:5432,bind=127.0.0.1,fork,reuseaddr TCP:midja.org:5432
