# services.d/rserve/run
#!/usr/bin/with-contenv bash
exec 2>&1
cd /usr/src/app
exec R CMD Rserve --quiet --vanilla
